# home/.config/systemd/user/llama-server.service
[Unit]
Description=Llama Server
After=network.target

[Service]
OOMScoreAdjust=500
# Command to start the server
ExecStart=/usr/bin/llama-server \
  --models-preset /home/sidney/models/config.ini \
  --host 0.0.0.0 \
  --port 11434 \
  --models-max 1

# Restart behavior (optional)
Restart=always
RestartSec=5

# Optional: Add environment variables if needed
#Environment="ENV_NAME=value"

[Install]
WantedBy=multi-user.target
